# deploy/wecare-parceiros-frontend.service
[Unit]
Description=WeCare Parceiros — Frontend Next.js
After=network.target

[Service]
Type=simple
User=wecare
WorkingDirectory=/opt/wecare-parceiros/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
